# bin/ubuntu.service
[Unit]
Description=DCapture H2 Console
After=network.target

[Service]
Type=forking

User=tomcat
Group=tomcat

Environment="JAVA_HOME=/usr/lib/jvm/java-11-openjdk-amd64"
Environment="JAVA_OPTS=-Djava.security.egd=file:///dev/urandom -Djava.awt.headless=true"

ExecStart=/opt/h2/bin/startup.sh
ExecStop=/opt/h2/bin/shutdown.sh

[Install]
WantedBy=multi-user.target